# roles/UBUNTU22-CIS/templates/etc/systemd/system/aidecheck.service.j2
[Unit]
Description=Aide Check

[Service]
Type=simple
ExecStart={{ ubtu22cis_aide_systemd_exec }}

[Install]
WantedBy=multi-user.target
